# xivdb-v3-sysops/crontab/SyncCharacterUpdate2
#
# SyncCharacterUpdate2 (21,22,23,24)
#

# Character updates
* * * * * /usr/bin/php /home/xivdb/xivsync/bin/console app:auto characters_update 21-Characters > /home/xivdb/logs/characters_update_21.log
* * * * * /usr/bin/php /home/xivdb/xivsync/bin/console app:auto characters_update 22-Characters > /home/xivdb/logs/characters_update_22.log
* * * * * /usr/bin/php /home/xivdb/xivsync/bin/console app:auto characters_update 23-Characters > /home/xivdb/logs/characters_update_23.log
* * * * * /usr/bin/php /home/xivdb/xivsync/bin/console app:auto characters_update 24-Characters > /home/xivdb/logs/characters_update_24.log

# Friends update
* * * * * /usr/bin/php /home/xivdb/xivsync/bin/console app:auto characters_friends 21-Characters > /home/xivdb/logs/characters_friends_21.log
* * * * * /usr/bin/php /home/xivdb/xivsync/bin/console app:auto characters_friends 22-Characters > /home/xivdb/logs/characters_friends_22.log
* * * * * /usr/bin/php /home/xivdb/xivsync/bin/console app:auto characters_friends 23-Characters > /home/xivdb/logs/characters_friends_23.log
* * * * * /usr/bin/php /home/xivdb/xivsync/bin/console app:auto characters_friends 24-Characters > /home/xivdb/logs/characters_friends_24.log

# Run Logs-Manager
0 * * * * /usr/bin/php /home/xivdb/xivsync/bin/console app:logs > /home/xivdb/logs/logs_manager.log
