# xivdb-v3-sysops/crontab/SyncAchievementUpdate2
#
# SyncAchievementUpdate2 (21,22,23,24)
#

# Achievement updates
* * * * * /usr/bin/php /home/xivdb/xivsync/bin/console app:auto characters_achievements 21-Achievements > /home/xivdb/logs/characters_achievements_21.log
* * * * * /usr/bin/php /home/xivdb/xivsync/bin/console app:auto characters_achievements 22-Achievements > /home/xivdb/logs/characters_achievements_22.log
* * * * * /usr/bin/php /home/xivdb/xivsync/bin/console app:auto characters_achievements 23-Achievements > /home/xivdb/logs/characters_achievements_23.log
* * * * * /usr/bin/php /home/xivdb/xivsync/bin/console app:auto characters_achievements 24-Achievements > /home/xivdb/logs/characters_achievements_24.log

# Run Logs-Manager
0 * * * * /usr/bin/php /home/xivdb/xivsync/bin/console app:logs > /home/xivdb/logs/logs_manager.log
